# deploy/systemd/sanyang-customer-order.service
[Unit]
Description=Sanyang Customer Order Admin (gunicorn :3003)
After=network.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/www/feijihe/stable
Environment=PYTHONUNBUFFERED=1
ExecStart=/www/feijihe/stable/venv/bin/gunicorn -w 2 -b 0.0.0.0:3003 --timeout 120 app_customer_order:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
