# deploy/systemd/sanyang-cs.service
[Unit]
Description=Sanyang CS (gunicorn :3001)
After=network.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/www/feijihe/stable
Environment=PYTHONUNBUFFERED=1
ExecStart=/www/feijihe/stable/venv/bin/gunicorn -w 2 -b 0.0.0.0:3001 --timeout 120 app_cs:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
