# deploy/systemd/sanyang-cs.service
[Unit]
Description=Sanyang CS (app_cs.py :3001)
After=network.target mysql.service

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/www/feijihe/stable
Environment=PYTHONUNBUFFERED=1
ExecStart=/www/feijihe/stable/venv/bin/python3 app_cs.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
